# settings/app/etc/supervisord.d/laravel-worker.ini
[program:laravel-worker-notification]
process_name=%(program_name)s_%(process_num)02d
command=/usr/bin/php %(ENV_APP_DIRECTORY_PATH)s/artisan queue:work sqs-notification --sleep=3 --tries=3 --max-time=3600
autostart=true
autorestart=true
numprocs=1
redirect_stderr=true
stdout_logfile=%(ENV_EFS_APPLICATION_MOUNT_PATH)s/logs/supervisor/laravel-worker-notification.log
stopwaitsecs=3600

[program:laravel-worker-email]
process_name=%(program_name)s_%(process_num)02d
command=/usr/bin/php %(ENV_APP_DIRECTORY_PATH)s/artisan queue:work sqs-email --sleep=3 --tries=3 --max-time=3600
autostart=true
autorestart=true
numprocs=1
redirect_stderr=true
stdout_logfile=%(ENV_EFS_APPLICATION_MOUNT_PATH)s/logs/supervisor/laravel-worker-email.log
stopwaitsecs=3600
